# crontab
0 */3 * * * cd .. && cd .. && cd ecomm_comp && python price_gathering_main.py